# python/1.3.2/test/esl-working
#/bin/sh

# chkconfig: 2345 99 70
# description: Starts and Stops ESL system.

case "$1" in
start)
	python /opt/esl-working/main.py start
;;
stop)
	python /opt/esl-working/main.py stop
;;
restart)
	python /opt/esl-working/main.py restart
;;
*)
	echo "Usage: esl-working {start|stop|restart}"
	exit 1
esac
exit 0
